# infra/backup/root.crontab
0 2 * * * /backup/backup.sh >> /var/backups/backup.log 2>&1
